# docker-build/conf.d/cron.d/cron-each-minute
#
#  Each minute
#
# m h dom mon dow user  command
* * * * * root /usr/local/bin/cron-each-minute.sh </dev/null &>/dev/null &
